# systemd/subrosa.service
[Unit]
Description=Subrosa AI Chief of Staff
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ati/Dev/subrosa
EnvironmentFile=/home/ati/.subrosa/.env
ExecStart=/home/ati/Dev/subrosa/.venv/bin/python -m subrosa
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
